# rootfs/etc/services.d/xorg/run
#!/usr/bin/with-contenv
Xorg\
    +extension RANDR\
    +extension RENDER\
    -config /etc/xpra/xorg.conf\
    -dpi $XORG_DPI\
    -logfile /tmp/Xorg.log\
    -nolisten tcp\
    -noreset\
    -xinerama